# tools/docker/tuleap-aio-c7/supervisor.d/backend_workers.ini
[program:backend_workers]
command = /usr/bin/tuleap worker:supervisor start
user=codendiadm
group=codendiadm
startsecs = 0
