# tools/docker/tuleap-aio-c7/supervisor.d/fpm.ini
[program:fpm]
command=/opt/remi/php73/root/usr/sbin/php-fpm --nodaemonize
redirect_stderr=true
process_name = php-fpm
autorestart = true
