# tools/docker/tuleap-aio-c7/supervisor.d/fpm.ini
[program:fpm]
command=/usr/share/tuleap/tools/docker/tuleap-aio-c7/supervisor.d/start-tuleap-php-fpm.sh
redirect_stderr=true
process_name = php-fpm
autorestart = true
