# tools/docker/tuleap-aio-c7/supervisor.d/mediawiki_fpm.ini
[program:mediawiki_fpm]
command=/opt/remi/php74/root/usr/sbin/php-fpm --nodaemonize --fpm-config /usr/share/tuleap/plugins/mediawiki_standalone/etc/php-fpm/mediawiki-tuleap.conf
redirect_stderr=true
process_name = php-fpm
autorestart = true
