# tools/docker/tuleap-aio-c7/supervisor.d/mysqld.ini
[program:mysql]
command = /opt/rh/rh-mysql57/root/usr/libexec/mysqld-scl-helper enable rh-mysql57 -- /opt/rh/rh-mysql57/root/usr/libexec/mysqld --basedir=/opt/rh/rh-mysql57/root/usr --pid-file=/var/run/rh-mysql57-mysqld/mysqld.pid
startsecs = 5
stopwaitsecs=300
user=mysql
group=mysql
