# tools/docker/tuleap-aio-c7/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true
process_name = nginx
autorestart = false
