# tools/docker/tuleap-aio-c7/supervisor.d/postfix.ini
[program:postfix]
process_name = master
command = /usr/sbin/postfix -c /etc/postfix start
startsecs = 0
autorestart = false
