# tools/docker/tuleap-aio-c7/supervisor.d/postfix.ini
[program:postfix]
process_name = master
command = /usr/sbin/postfix -D -c /etc/postfix start
autorestart = false
