# tools/docker/tuleap-aio-c7/supervisor.d/realtime.ini
[program:tuleap-realtime]
command=/usr/share/tuleap/tools/docker/tuleap-aio-c7/supervisor.d/start-tuleap-realtime.sh
process_name=tuleap-realtime
redirect_stderr=true
autorestart=true
