# tools/docker/tuleap-aio-c7/supervisor.d/rsyslog.ini
[program:rsyslog]
command = rsyslogd -n -f /etc/rsyslog.conf
startsecs = 5
stopwaitsecs = 5
