# tools/docker/tuleap-aio-c7/supervisor.d/smokescreen.ini
[program:tuleap-smokescreen]
command=/usr/share/tuleap/tools/docker/tuleap-aio-c7/supervisor.d/start-tuleap-smokescreen.sh
process_name=tuleap-smokescreen
redirect_stderr=true
autorestart=true
