# tools/tuleap-docker-supervisor.d/crond.ini
[program:crond]
command=/usr/sbin/crond -s -n