# tools/tuleap-docker-supervisor.d/fpm.ini
[program:fpm]
command=/usr/share/tuleap/tools/tuleap-docker-supervisor.d/start-tuleap-php-fpm.sh
redirect_stderr=true
process_name = php-fpm
autorestart = true
