# tools/tuleap-docker-supervisor.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
redirect_stderr=true
process_name = httpd
