# tools/tuleap-docker-supervisor.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -E /dev/stderr -D
redirect_stderr=true
