# src/additional-packages/tuleap-mercure/tuleap-mercure-config-change.path
[Path]
PathChanged=/etc/tuleap/conf/mercure.env

[Install]
WantedBy=multi-user.target
