# plugins/fts_meilisearch/additional-packages/tuleap-meilisearch-server/tuleap-meilisearch-config-change.path
[Path]
PathChanged=/var/lib/tuleap/fts_meilisearch_server/meilisearch-master-key.env

[Install]
WantedBy=multi-user.target
