# plugins/fts_meilisearch/additional-packages/tuleap-meilisearch-server/tuleap-meilisearch-config-change.service
[Unit]
Description=Restart Meilisearch server for Tuleap on configuration change
Documentation=https://docs.tuleap.org

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart tuleap-meilisearch.service
