# src/tuleap-cfg/resources/realtime.ini
[program:tuleap-realtime]
command=/usr/bin/node /usr/lib/node_modules/tuleap-realtime/server.js --config=/etc/tuleap-realtime/config.json
redirect_stderr= true
autorestart = true
