# src/tuleap-cfg/resources/realtime.ini
[program:tuleap-realtime]
command=/usr/lib/tuleap-realtime/tuleap-realtime --config=/etc/tuleap-realtime/config.json
redirect_stderr= true
autorestart = true
